# scripts/minimax-deploy/com.headroom.minimax-token-refresher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.headroom.minimax-token-refresher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/axel/.mavis/bin/minimax-token-refresher.sh</string>
  </array>
  <!-- Ogni 6 ore. Lo script è idempotente: se il token non è cambiato, esce subito. -->
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/axel/.headroom/logs/token-refresher.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/axel/.headroom/logs/token-refresher.log</string>
  <key>ExitTimeOut</key>
  <integer>15</integer>
</dict>
</plist>